# deployment/tethys_deploy/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=service nginx start
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log